# daemons/raise-awareness.service
[Unit]
Description=Raise Awareness Daemon - Security anomaly monitor
After=network.target
Wants=entropy-economy.service factory.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace
ExecStart=/root/.openclaw/workspace/daemons/raise-awareness.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raise-awareness

# Security
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=read-only  # disabled due to log writing
PrivateTmp=true

[Install]
WantedBy=multi-user.target